# systemd/anpr_core.service
[Unit]
Description=ANPR Gate System — Core Recognition Service
After=network.target

[Service]
Type=simple
User=machinemustlearn
WorkingDirectory=/home/machinemustlearn/WORKSPACE/ANPR_System
# Use the virtual‑environment python interpreter
ExecStart=/home/machinemustlearn/WORKSPACE/ANPR_System/venv/bin/python -m src.main --config /home/machinemustlearn/WORKSPACE/ANPR_System/config.yaml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
