# systemd/anpr_web.service
[Unit]
Description=ANPR Gate System — Web Dashboard
After=network.target

[Service]
Type=simple
User=machinemustlearn
WorkingDirectory=/home/machinemustlearn/WORKSPACE/ANPR_System
ExecStart=/home/machinemustlearn/WORKSPACE/ANPR_System/venv/bin/python -m web.app
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
